# backend/claude-panel.service
# systemd unit -> /etc/systemd/system/claude-panel.service
#   sudo cp backend/claude-panel.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now claude-panel
#   journalctl -u claude-panel -f      # contraseña + inscripción 2FA al arrancar
[Unit]
Description=Panel web Claude (agente sobre esta VPS)
After=network-online.target
Wants=network-online.target

[Service]
User=adelpozor
WorkingDirectory=/home/adelpozor/claude-panel/backend
# El CLI `claude` vive en ~/.local/bin, que no está en el PATH por defecto de systemd.
Environment=PATH=/home/adelpozor/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/bin:/snap/bin
ExecStart=/usr/bin/python3 /home/adelpozor/claude-panel/backend/server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
